# deploy/sub2api-subsite-agent.service
[Unit]
Description=Sub2API Subsite Agent - Lightweight Edge Executor
Documentation=https://github.com/Wei-Shaw/sub2api
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sub2api
Group=sub2api
WorkingDirectory=/opt/sub2api-subsite-agent
EnvironmentFile=/etc/sub2api/subsite-agent.env
ExecStart=/opt/sub2api-subsite-agent/sub2api-subsite-agent
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sub2api-subsite-agent

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/sub2api-subsite-agent

[Install]
WantedBy=multi-user.target
